# pueue/pueued.service
[Unit]
Description=Pueue Daemon - CLI process scheduler and manager

[Service]
Restart=no
ExecStart=/home/okr/.cargo/bin/pueued -vvv

[Install]
WantedBy=default.target
